# init.conf
description "Mesos-DNS"

start on runlevel [2345]
stop on runlevel [!2345]

console log

respawn
respawn limit 10 5

exec /usr/sbin/mesos-dns -config /etc/mesos/mesos-dns.json -logtostderr
